# services/knmi-data-refresh.service
[Unit]
Description=Rebuild KNMI dashboard data (dashboard/data/*.json)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=sausje
WorkingDirectory=/home/sausje/personal/knmi
ExecStart=/home/sausje/personal/knmi/.venv/bin/python scripts/build_data.py
